# Persistence/Potential_Persistence_Via_COM_Hijacking_From_Suspicious_Locations.kql
// Author: Nasreddine Bencherchali (Nextron Systems)
// Date: 2022/07/28
// Level: high
// Description: Detects potential COM object hijacking where the "Server" (In/Out) is pointing to a supsicious or unsuale location
// Tags: attack.persistence, attack.t1546.015
DeviceRegistryEvents
| where (RegistryValueData contains "\\Users\\Public\\" or RegistryValueData contains "\\Desktop\\" or RegistryValueData contains "\\AppData\\Local\\Temp\\" or RegistryValueData contains "\\Downloads\\" or RegistryValueData contains "\\Windows\\Temp\\" or RegistryValueData contains "\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\" or RegistryValueData contains "\\System32\\spool\\drivers\\color\\" or RegistryValueData contains "%temp%" or RegistryValueData contains "%tmp%" or RegistryValueData contains "%appdata%") and ActionType =~ "RegistryValueSet" and (RegistryKey endswith "\\InprocServer32\\(Default)" or RegistryKey endswith "\\LocalServer32\\(Default)") and (RegistryKey startswith "HKEY_LOCAL_MACHINE\\CLASSES\\CLSID" or RegistryKey startswith "HKEY_CLASSES_ROOT\\CLSID" or RegistryKey startswith "HKCU\\Software\\Classes\\CLSID" or RegistryKey startswith "HKEY_CURRENT_USER\\Software\\Classes\\CLSID")