# KQL/Leviathan_Registry_Key_Activity.kql
// Author: Aidan Bracher
// Date: 2020/07/07
// Level: critical
// Description: Detects registry key used by Leviathan APT in Malaysian focused campaign
// Tags: attack.persistence, attack.t1547.001
DeviceRegistryEvents
| where RegistryKey =~ "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\ntkd"