# KQL/Suspicious_Command_Patterns_In_Scheduled_Task_Creation.kql
// Author: Florian Roth (Nextron Systems)
// Date: 2022/02/23
// Level: high
// Description: Detects scheduled task creation using "schtasks" that contain potentially suspicious or uncommon commands
// Tags: attack.execution, attack.t1053.005
DeviceProcessEvents
| where (ProcessCommandLine contains "/Create " and FolderPath endswith "\\schtasks.exe") and (((ProcessCommandLine contains "/sc minute " or ProcessCommandLine contains "/ru system ") and (ProcessCommandLine contains "cmd /c" or ProcessCommandLine contains "cmd /k" or ProcessCommandLine contains "cmd /r" or ProcessCommandLine contains "cmd.exe /c " or ProcessCommandLine contains "cmd.exe /k " or ProcessCommandLine contains "cmd.exe /r ")) or (ProcessCommandLine contains " -decode " or ProcessCommandLine contains " -enc " or ProcessCommandLine contains " -w hidden " or ProcessCommandLine contains " bypass " or ProcessCommandLine contains " IEX" or ProcessCommandLine contains ".DownloadData" or ProcessCommandLine contains ".DownloadFile" or ProcessCommandLine contains ".DownloadString" or ProcessCommandLine contains "/c start /min " or ProcessCommandLine contains "FromBase64String" or ProcessCommandLine contains "mshta http" or ProcessCommandLine contains "mshta.exe http") or ((ProcessCommandLine contains "\\AppData\\" and ProcessCommandLine contains "%AppData%" and ProcessCommandLine contains "%Temp%" and ProcessCommandLine contains "%tmp%" and ProcessCommandLine contains "C:\\Windows\\Temp\\") and (ProcessCommandLine contains "/xml C:\\Users\\" and ProcessCommandLine contains "cscript" and ProcessCommandLine contains "curl" and ProcessCommandLine contains "wscript")))