# KQL/Suspicious_Powershell_In_Registry_Run_Keys.kql
// Author: frack113, Florian Roth
// Date: 2022/03/17
// Level: medium
// Description: Detects potential PowerShell commands or code within registry run keys
// Tags: attack.persistence, attack.t1547.001
DeviceRegistryEvents
| where (RegistryValueData contains "powershell" or RegistryValueData contains "pwsh " or RegistryValueData contains "FromBase64String" or RegistryValueData contains ".DownloadFile(" or RegistryValueData contains ".DownloadString(" or RegistryValueData contains " -w hidden " or RegistryValueData contains " -w 1 " or RegistryValueData contains "-windowstyle hidden" or RegistryValueData contains "-window hidden" or RegistryValueData contains " -nop " or RegistryValueData contains " -encodedcommand " or RegistryValueData contains "-ExecutionPolicy Bypass" or RegistryValueData contains "Invoke-Expression" or RegistryValueData contains "IEX (" or RegistryValueData contains "Invoke-Command" or RegistryValueData contains "ICM -" or RegistryValueData contains "Invoke-WebRequest" or RegistryValueData contains "IWR " or RegistryValueData contains " -noni " or RegistryValueData contains " -noninteractive ") and ActionType =~ "RegistryValueSet" and RegistryKey contains "\\Software\\Microsoft\\Windows\\CurrentVersion\\Run"