# Execution/Suspicious_Schtasks_From_Env_Var_Folder.kql
// Author: Florian Roth (Nextron Systems)
// Date: 2022/02/21
// Level: medium
// Description: Detects Schtask creations that point to a suspicious folder or an environment variable often used by malware
// Tags: attack.execution, attack.t1053.005
DeviceProcessEvents
| where (((ProcessCommandLine contains ":\\Perflogs" or ProcessCommandLine contains ":\\Windows\\Temp" or ProcessCommandLine contains "\\AppData\\Local\\" or ProcessCommandLine contains "\\AppData\\Roaming\\" or ProcessCommandLine contains "\\Users\\Public" or ProcessCommandLine contains "%AppData%" or ProcessCommandLine contains "%Public%") and (ProcessCommandLine contains " /create " and FolderPath endswith "\\schtasks.exe")) or (InitiatingProcessCommandLine endswith "\\svchost.exe -k netsvcs -p -s Schedule" and (ProcessCommandLine contains ":\\Perflogs" or ProcessCommandLine contains ":\\Windows\\Temp" or ProcessCommandLine contains "\\Users\\Public" or ProcessCommandLine contains "%Public%"))) and (not(((ProcessCommandLine contains "/Create /Xml \"C:\\Users\\" and ProcessCommandLine contains "\\AppData\\Local\\Temp\\.CR." and ProcessCommandLine contains "Avira_Security_Installation.xml") or ((ProcessCommandLine contains ".tmp\\UpdateFallbackTask.xml" or ProcessCommandLine contains ".tmp\\WatchdogServiceControlManagerTimeout.xml" or ProcessCommandLine contains ".tmp\\SystrayAutostart.xml" or ProcessCommandLine contains ".tmp\\MaintenanceTask.xml") and (ProcessCommandLine contains "/Create /F /TN" and ProcessCommandLine contains "/Xml " and ProcessCommandLine contains "\\AppData\\Local\\Temp\\is-" and ProcessCommandLine contains "Avira_")) or (ProcessCommandLine contains "\\AppData\\Local\\Temp\\" and ProcessCommandLine contains "/Create /TN \"klcp_update\" /XML " and ProcessCommandLine contains "\\klcp_update_task.xml") or ((ProcessCommandLine contains "update_task.xml" or ProcessCommandLine contains "/Create /TN TVInstallRestore /TR") or InitiatingProcessCommandLine contains "unattended.ini"))))